# Exc8/logview.service
[Unit]
Description=My service for view log

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/logview
ExecStart=/opt/logview.sh $WORD $LOG
